# systemd/audio.slice

[Unit]
Description=Audio processes
Before=slices.target

[Slice]
MemoryAccounting=true
IOAccounting=true
CPUAccounting=true

