# systemd/websocket-ready.target
[Unit]
Description=WebSocket Ready (CamillaDSP)
Requires=websocket-ready.service
After=websocket-ready.service

